# cookbooks/rails/remote_files/etc/profile.d/rails.sh
export RAILS_DATABASE_DEV_USER=rooster
export RAILS_DATABASE_DEV_HOST=192.168.1.40
export RAILS_DATABASE_DEV_PORT=5432
export RAILS_DATABASE_DEV_PASSWORD=G8BrneBL

export SECRET_KEY_BASE=c982fd11933ba078cf2bf4fecb62a38e1d3c6b95852a0c83fe1652786c9509f7f59cece77bc556581968aa647e2aa282d8239620e2ac2f96298a8b53e4a67387
